# root/ryzen/systemd/system/lvm2-monitor.service.d/override.conf
[Unit]
Requires=dm-event.socket tailscaled.service snap.microceph.daemon.service
After=dm-event.socket dm-event.service tailscaled.service snap.microceph.daemon.service
