# root/zen/systemd/system/lvm2-monitor.service.d/override.conf
[Unit]
Requires=dm-event.socket snap.microceph.daemon.service
After=dm-event.socket dm-event.service snap.microceph.daemon.service
