# root/ryzen/systemd/system/snap.microceph.daemon.service.d/override.conf
[Unit]
Requires=snap-microceph-707.mount tailscaled.service
After=snap-microceph-707.mount network.target snapd.apparmor.service tailscaled.service
